# home/private_dot_config/systemd/user/app-com.nextcloud.desktopclient.nextcloud@autostart.service.d/override.conf
[Unit]
Requires=keepassxc.service
Wants=tray.service 
After=tray.service keepassxc.service graphical-session.target
PartOf=

[Service]
Slice=Background.slice
Restart=always
RestartSec=20s
