# private_dot_config/systemd/user/emacs.service.d/override.conf
[Service]
ExecStartPre=/usr/bin/sleep 10

[Install]
WantedBy=wm.target
