# private_dot_config/systemd/user/app-Flameshot@autostart.service.d/override.conf
[Unit]
Wants=tray.service
After=tray.service graphical-session.target

[Service]
Restart=always
RestartSec=30s
