# home/private_dot_config/systemd/user/app-org.fcitx.Fcitx5@autostart.service.d/override.conf
[Unit]
Wants=tray.service
After=tray.service graphical-session.target

[Service]
Restart=always
RestartSec=20s
KillMode=control-group