# private_dot_config/systemd/user/duplicati.service.d/override.conf
[Unit]
Wants=tray.service
Requires=wait-for-dns.service
After=tray.service wait-for-dns.service

[Service]
Slice=background.slice
Restart=always
RestartSec=20s
