# private_dot_config/systemd/user/autolock.timer
[Unit]
Description=Automatically lock system between 17-17:15

[Timer]
OnCalendar=17:00..15:00

[Install]
WantedBy=timers.target
