# private_dot_config/systemd/user/xdg-autostarts.service
[Unit]
Description=Bring in XDG-Autostarts auto-generated by systemd
Wants=xdg-desktop-autostart.target

[Service]
Type=oneshot
ExecStart=/usr/bin/sleep 1

[Install]
WantedBy=wm.target
