# home/private_dot_config/systemd/user/autolock.service
[Unit]
Description=Lock computer, checking if it's already locked
ConditionFileIsExecutable=%h/.local/bin/check_i3lock

[Service]
Type=oneshot
ExecStart=%h/.local/bin/check_i3lock
