# home/private_dot_config/systemd/user/aw-server.service
[Unit]
Description=Activity Watch server

[Service]
Type=exec
Restart=on-failure
RestartSec=20s
ExecStart=/usr/bin/aw-server

[Install]
WantedBy=wm.target
