# home/private_dot_config/systemd/user/copyq.service
[Unit]
Description=Copyq clipboard manager
Documentation=man:copyq(1)
Wants=tray.service
After=tray.service
ConditionFileIsExecutable=/usr/bin/copyq

[Service]
Type=exec
RestartSec=20s
Restart=always
ExecStart=copyq

[Install]
WantedBy=wm.target
