# home/private_dot_config/systemd/user/xdg-autostarts.service
[Unit]
Description=Bring in XDG-Autostarts auto-generated by systemd
Wants=xdg-desktop-autostart.target
ConditionFileIsExecutable=/usr/bin/sleep

[Service]
Type=oneshot
ExecStart=sleep 1

[Install]
WantedBy=wm.target
