# private_dot_config/systemd/user/app-com.nextcloud.desktopclient.nextcloud@autostart.service.d/override.conf
[Unit]
Wants=tray.service keepassxc.service
After=tray.service keepassxc.service graphical-session.target

[Service]
Slice=Background.slice
Restart=always
RestartSec=20s
